# tools/docker/build/bin/supervisor.d/10-vnc.conf
[program:vnc]
command=/usr/local/bin/start_vnc.sh
user=vimacc
directory=/home/vimacc
environment=HOME="/home/vimacc",USER="vimacc",LOGNAME="vimacc",SHELL="/bin/bash",RESOLUTION="%(ENV_RESOLUTION)s"
autostart=true
autorestart=unexpected
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=10
