# tools/docker/build/bin/supervisor.d/20-vimacc-core.conf
; 20-vimacc-core.conf
; vimacc core services managed by supervisord


; ===== BASIC services (always) =====
[program:AccVimaccRoot]
command=/opt/Accellence/vimacc/bin/AccVimaccRoot -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=20
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccReportClient]
command=/opt/Accellence/vimacc/bin/AccVimaccReportClient -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=20
stopsignal=TERM
stopwaitsecs=90


; ===== Management servers (central Video Config Servers) =====
[program:AccVimaccConfig]
command=/opt/Accellence/vimacc/bin/AccVimaccConfig -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=30
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccReportServer]
command=/opt/Accellence/vimacc/bin/AccVimaccReportServer -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=30
stopsignal=TERM
stopwaitsecs=120


; ===== Sub-Nodes (statt obiger Mgmt-Server) =====
[program:AccVimaccConfigSlave]
command=/opt/Accellence/vimacc/bin/AccVimaccConfigSlave -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=30
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccConfigProxy]
command=/opt/Accellence/vimacc/bin/AccVimaccConfigProxy -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=30
stopsignal=TERM
stopwaitsecs=90

; ===== Recommended general services =====
[program:AccVimaccSystemMonitor]
command=/opt/Accellence/vimacc/bin/AccVimaccSystemMonitor -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=40
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccEventManager]
command=/opt/Accellence/vimacc/bin/AccVimaccEventManager -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=40
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccScheduler]
command=/opt/Accellence/vimacc/bin/AccVimaccScheduler -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=40
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccAssembliesManager]
command=/opt/Accellence/vimacc/bin/AccVimaccAssembliesManager -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=40
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccControlInterface]
command=/opt/Accellence/vimacc/bin/AccVimaccControlInterface -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=40
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccSnmpAgent]
command=/opt/Accellence/vimacc/bin/AccVimaccSnmpAgent -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=40
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccEmailSender]
command=/opt/Accellence/vimacc/bin/AccVimaccEmailSender -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=40
stopsignal=TERM
stopwaitsecs=90


; ===== Dongle-Manager =====
[program:AccVimaccDongle]
command=/opt/Accellence/vimacc/bin/AccVimaccDongle -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=45
stopsignal=TERM
stopwaitsecs=90


; ===== HID-Controller =====
[program:AccVimaccHIDController]
command=/opt/Accellence/vimacc/bin/AccVimaccHIDController -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=45
stopsignal=TERM
stopwaitsecs=90


; ===== Device handling / streaming / recording / playback / export =====
[program:AccVimaccInterface]
command=/opt/Accellence/vimacc/bin/AccVimaccInterface -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=50
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccServer]
command=/opt/Accellence/vimacc/bin/AccVimaccServer -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=50
stopsignal=TERM
stopwaitsecs=600

[program:AccVimaccExport]
command=/opt/Accellence/vimacc/bin/AccVimaccExport -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=50
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccIpAlarmReceiver]
command=/opt/Accellence/vimacc/bin/AccVimaccIpAlarmReceiver -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=50
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccDigIOInterface]
command=/opt/Accellence/vimacc/bin/AccVimaccDigIOInterface -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=50
stopsignal=TERM
stopwaitsecs=90


; ===== Proxies (zentral empfohlen) =====
[program:AccVimaccInterfaceProxy]
command=/opt/Accellence/vimacc/bin/AccVimaccInterfaceProxy -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=60
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccPlaybackProxy]
command=/opt/Accellence/vimacc/bin/AccVimaccPlaybackProxy -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=60
stopsignal=TERM
stopwaitsecs=90


; ===== Recording-Sync =====
[program:AccVimaccServerSync]
command=/opt/Accellence/vimacc/bin/AccVimaccServerSync -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=65
stopsignal=TERM
stopwaitsecs=90


; ===== System integration =====
[program:AccVimaccRTSPServer]
command=/opt/Accellence/vimacc/bin/AccVimaccRTSPServer -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=70
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccHttpServer]
command=/opt/Accellence/vimacc/bin/AccVimaccHttpServer -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=70
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccDisplayInterface]
command=/opt/Accellence/vimacc/bin/AccVimaccDisplayInterface -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=70
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccFtpAlarmReceiver]
command=/opt/Accellence/vimacc/bin/AccVimaccFtpAlarmReceiver -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=70
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccFtpUploader]
command=/opt/Accellence/vimacc/bin/AccVimaccFtpUploader -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=70
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccSmsGateway]
command=/opt/Accellence/vimacc/bin/AccVimaccSmsGateway -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=70
stopsignal=TERM
stopwaitsecs=90


; ===== Analytics =====
[program:AccVimaccMotionDetector]
command=/opt/Accellence/vimacc/bin/AccVimaccMotionDetector -e
user=vimacc
autostart=true
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=80
stopsignal=TERM
stopwaitsecs=90

[program:AccVimaccLicensePlateDetector]
command=/opt/Accellence/vimacc/bin/AccVimaccLicensePlateDetector -e
user=vimacc
autostart=false
autorestart=true
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
priority=80
stopsignal=TERM
stopwaitsecs=90

